# src/plugins/analysis/init_systems/test/data/etc/init/baz.conf
description "Simple Baz application"
author      "Foo"

start on filesystem or runlevel [2345]
stop on shutdown

script
    export HOME="/home/foo"
    echo $$ > /var/run/baz.pid
    exec /usr/local/bin/baz /home/foo/config.ru
end script

pre-start script
    echo "[`date`] baz starting..." >> /var/log/baz.log
end script

pre-stop script
    rm /var/run/baz.pid
    echo "[`date`] baz stopping..." >> /var/log/baz.log
end script

exec /bin/baz.sh -runonce \
-silent