# fact/plugins/analysis/init_systems/test/data/etc/sv/example/run
#!/bin/sh -e
exec 2>&1
exec chpst -u foo /opt/example/foo-service.sh