# fact/plugins/analysis/init_systems/test/data/etc/service/lighttpd/run
#!/bin/sh -e
sv -w7 check postgresql
exec 2>&1 myprocess \
last line